# systemd/smartgrill.service.template
[Unit]
Description=SmartGrill AT-02 webserver
After=network-online.target bluetooth.service
Wants=network-online.target
Requires=bluetooth.service

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__DIR__
EnvironmentFile=__DIR__/.env
ExecStart=__DIR__/.venv/bin/uvicorn smartgrill.app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
